# deploy/account-monitor.service
[Unit]
Description=theta account monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/theta
Environment=HOME=/root
Environment=RUST_LOG=info
ExecStart=/root/theta/target/release/theta ops account-monitor --loop --every-seconds 300 --market-hours-only true --account firstrade
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
