# deploy/capture-signals.service
[Unit]
Description=theta capture-signals sampler
After=network-online.target theta-daemon@%i.service
Wants=network-online.target
Requires=theta-daemon@%i.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h
Environment=HOME=%h
Environment=RUST_LOG=info
EnvironmentFile=-%h/.config/theta/config.env
EnvironmentFile=-%h/.config/theta/capture-signals.env
ExecStartPre=/bin/sh -c 'socket_path="${THETA_SOCKET_PATH:-$HOME/.theta/run/theta.sock}"; for _ in $(seq 1 30); do [ -S "$socket_path" ] && exit 0; sleep 1; done; echo "theta-daemon socket not ready: $socket_path" >&2; exit 1'
ExecStart=/usr/local/bin/theta signals capture --loop --every-seconds 300 --market-hours-only
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
