# deploy/capture-signals.service
[Unit]
Description=theta capture-signals sampler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/theta
Environment=HOME=%h
Environment=RUST_LOG=info
ExecStart=%h/theta/target/release/theta signals capture --loop --every-seconds 300
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
