# deploy/capture-signals.service
[Unit]
Description=theta capture-signals sampler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/theta
Environment=HOME=%h
Environment=RUST_LOG=info
EnvironmentFile=-%h/.config/theta/capture-signals.env
ExecStart=%h/theta/target/release/capture-signals --loop --every-seconds 300 --market-hours-only
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
